# app/websockets.conf
[program:websockets]
command=/usr/bin/php /home/laravel-echo/laravel-websockets/artisan websockets:serve --port=6002
numprocs=1
autostart=true
autorestart=true
user=laravel-echo
